# etc/cron.d/sundowner-instagram
*/15 * * * * root mongo --eval "var DATABASE_NAME='sundowner_instagram';" /home/jhibberd/projects/sundowner/src/mongodb-js/update-score.js
