# etc/init/sundowner-www.conf

start on runlevel [2345]
stop on runlevel [06]

console log 
script
    exec python /opt/sundowner/www/main.py
end script

